# root/etc/services.d/zap/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid zap zap.sh -daemon -host 0.0.0.0 -newsession /tmp/zap/;
